# src/data_management/lawnberry-data.service
[Unit]
Description=LawnBerry Data Management Service
After=network-online.target lawnberry-communication.service
Wants=network-online.target
Requires=lawnberry-communication.service

[Service]
Type=simple
User=lawnberry
Group=lawnberry
WorkingDirectory=/var/lib/lawnberry
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.data_management.data_service
Restart=on-failure
RestartSec=3

# Hardening options
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
ReadOnlyPaths=/
ReadWritePaths=/var/lib/lawnberry /var/log/lawnberry /run /tmp
DeviceAllow=/dev/tty* rw
MemoryMax=512M

[Install]
WantedBy=multi-user.target
