# v2/lawnberry-v2/systemd/mower-core.service
[Unit]
Description=LawnBerry core controller
After=network-online.target
[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/lawnberry_pi/v2/lawnberry-v2
ExecStart=/usr/bin/env uv run python -m lawnberry.core
Restart=always
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/pi/lawnberry_pi/.env
[Install]
WantedBy=multi-user.target
