# src/weather/lawnberry-weather.service
[Unit]
Description=Lawnberry Weather Integration Service
Documentation=https://github.com/lawnberry/weather-integration
After=lawnberry-communication.service lawnberry-data.service
Requires=lawnberry-communication.service lawnberry-data.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.weather.weather_service
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=10
StartLimitBurst=2

# Resource limits
MemoryMax=256M
CPUQuota=50%
TasksMax=50

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Network access for API calls
PrivateNetwork=false

# Performance settings
Nice=5
IOSchedulingClass=2
IOSchedulingPriority=4

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-weather

# Watchdog removed (sd_notify not used by service)

[Install]
WantedBy=multi-user.target
