# lawnberry-rebuild/systemd/lawnberry-health.service
[Unit]
Description=LawnBerry Pi v2 System Health Monitor
Documentation=https://github.com/acredsfan/lawnberry_pi
After=lawnberry-backend.service
Wants=lawnberry-backend.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/lawnberry/lawnberry-rebuild
Environment=PYTHONPATH=/home/pi/lawnberry/lawnberry-rebuild

# Health monitoring script
ExecStart=/bin/bash -c 'while true; do \
    echo "$(date): Health check - Backend API, Database, Sensors" >> /home/pi/lawnberry/logs/health.log; \
    curl -f http://localhost:8081/health >/dev/null 2>&1 || echo "$(date): Backend API health check failed" >> /home/pi/lawnberry/logs/health.log; \
    sleep 30; \
done'

# Process management
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/lawnberry/logs

# Resource limits
LimitNOFILE=64
MemoryMax=32M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-health

[Install]
WantedBy=multi-user.target