# lawnberry-rebuild/systemd/lawnberry-acme-renew.timer
[Unit]
Description=LawnBerry Pi ACME Certificate Renewal Timer
Requires=lawnberry-acme-renew.service

[Timer]
# Run twice daily at 00:00 and 12:00
OnCalendar=00:00,12:00
# Run on boot after 10 minutes
OnBootSec=10min
# Add randomization to avoid load spikes
RandomizedDelaySec=30min
# Ensure persistent timer across reboots
Persistent=true

[Install]
WantedBy=timers.target