# src/vision/lawnberry-vision.service
[Unit]
Description=Lawnberry Computer Vision System
Documentation=https://github.com/lawnberry/computer-vision
After=lawnberry-hardware.service lawnberry-communication.service
Requires=lawnberry-hardware.service lawnberry-communication.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.vision.vision_manager
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=10
StartLimitBurst=2

# Resource limits
MemoryMax=1G
CPUQuota=150%
TasksMax=100

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry /dev
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Hardware access for camera and TPU
SupplementaryGroups=video

# Performance settings
Nice=-1
IOSchedulingClass=1
IOSchedulingPriority=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-vision

# Watchdog removed (sd_notify not used by service)

[Install]
WantedBy=multi-user.target
