# systemd/lawnberry-database.service
[Unit]
Description=LawnBerry Pi v2 Database Service
Documentation=https://github.com/acredsfan/lawnberry_pi
After=local-fs.target
Before=lawnberry-backend.service

[Service]
Type=oneshot
RemainAfterExit=true
User=pi
Group=pi

# Database initialization and migration
WorkingDirectory=/apps/lawnberry-pi
Environment=PYTHONPATH=/apps/lawnberry-pi
ExecStart=/bin/bash -lc 'mkdir -p ./data ./logs ./config && \
    ./.venv/bin/python ./scripts/init_database.py'

# Permissions and directory setup
ExecStartPre=/bin/mkdir -p /apps/lawnberry-pi/data
ExecStartPre=/bin/mkdir -p /apps/lawnberry-pi/logs
ExecStartPre=/bin/mkdir -p /apps/lawnberry-pi/config
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R pi:pi /apps/lawnberry-pi

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/apps/lawnberry-pi

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-database

[Install]
WantedBy=multi-user.target
