# systemd/lawnberry-cert-renewal.timer
[Unit]
Description=Daily LawnBerry Pi certificate renewal/validation timer
Documentation=README.md

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
