# systemd/lawnberry-backend.service
[Unit]
Description=LawnBerry Pi v2 Backend API Service
Documentation=https://github.com/acredsfan/lawnberry_pi
After=network-online.target
Wants=network-online.target
Requires=lawnberry-database.service

[Service]
Type=exec
User=pi
Group=pi
WorkingDirectory=/home/pi/lawnberry/lawnberry-rebuild
Environment=PYTHONPATH=/home/pi/lawnberry/lawnberry-rebuild
Environment=LAWNBERRY_ENV=production
Environment=LAWNBERRY_CONFIG_DIR=/home/pi/lawnberry/config
Environment=LAWNBERRY_LOG_DIR=/home/pi/lawnberry/logs
Environment=LAWNBERRY_DATA_DIR=/home/pi/lawnberry/data
Environment=SIM_MODE=0

# Main service command
ExecStart=/home/pi/lawnberry/lawnberry-rebuild/.venv/bin/uvicorn backend.src.main:app --host 0.0.0.0 --port 8081 --workers 1

# Health checks and monitoring
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Process management
Restart=always
RestartSec=10
RestartPreventExitStatus=2

# Security restrictions (ARM64/Pi OS specific)
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/lawnberry/logs /home/pi/lawnberry/data /home/pi/lawnberry/config
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Hardware access (required for GPIO, I2C, SPI)
SupplementaryGroups=gpio i2c spi dialout

# Resource limits
LimitNOFILE=4096
MemoryMax=512M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-backend

[Install]
WantedBy=multi-user.target