# systemd/lawnberry-camera.service
[Unit]
Description=LawnBerry Pi Camera Stream Service
After=network-online.target
Wants=network-online.target
PartOf=lawnberry-backend.service

[Service]
Type=simple
User=pi
Group=video
WorkingDirectory=/home/pi/lawnberry
Environment=SIM_MODE=0
ExecStart=/bin/bash -lc "PY='/home/pi/lawnberry/.venv/bin/python'; if [ ! -x \"$PY\" ]; then PY='python3'; fi; cd /home/pi/lawnberry && \"$PY\" -m backend.src.services.camera_stream_service"
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15

# Permissions for camera access
SupplementaryGroups=video render

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/lawnberry /var/lib/lawnberry /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-camera

[Install]
WantedBy=multi-user.target