# systemd/lawnberry-camera.service
[Unit]
Description=LawnBerry Pi Camera Stream Service
After=network-online.target
Wants=network-online.target
PartOf=lawnberry-backend.service

[Service]
Type=simple
User=pi
Group=video
WorkingDirectory=/home/pi/lawnberry
EnvironmentFile=-/home/pi/lawnberry/.env
Environment=SIM_MODE=0
Environment=LAWNBERRY_CAMERA_SOCKET=/run/lawnberry/camera.sock
RuntimeDirectory=lawnberry
RuntimeDirectoryMode=0755
# ExecStart-level assignments override EnvironmentFile values so production
# cannot accidentally embed a simulated camera owner or split socket paths.
ExecStart=/usr/bin/env SIM_MODE=0 LAWNBERRY_CAMERA_SOCKET=/run/lawnberry/camera.sock /home/pi/lawnberry/.venv/bin/python -m backend.src.services.camera_stream_service
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15

# Permissions for camera access
SupplementaryGroups=video render

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/lawnberry /var/lib/lawnberry /run/lawnberry

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-camera

[Install]
WantedBy=multi-user.target
