# systemd/lawnberry-acme-renew.service
[Unit]
Description=LawnBerry Pi ACME Certificate Renewal
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=.
ExecStart=/bin/bash -lc "PY='./.venv/bin/python'; if [ ! -x \"$PY\" ]; then PY='python3'; fi; cd . && \"$PY\" -m backend.src.cli.acme_renew"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-acme-renew

[Install]
WantedBy=multi-user.target