# systemd/lawnberry-database.service
[Unit]
Description=LawnBerry Pi v2 Database Service
Documentation=https://github.com/acredsfan/lawnberry_pi
After=local-fs.target
Before=lawnberry-backend.service

[Service]
Type=oneshot
RemainAfterExit=true
User=pi
Group=pi

# Database initialization and migration
WorkingDirectory=/home/pi/lawnberry
Environment=PYTHONPATH=/home/pi/lawnberry
ExecStart=/bin/bash -lc 'mkdir -p /home/pi/lawnberry/data /home/pi/lawnberry/logs /home/pi/lawnberry/config && \
    /home/pi/lawnberry/.venv/bin/python /home/pi/lawnberry/scripts/init_database.py'

# Permissions and directory setup
ExecStartPre=/bin/mkdir -p /home/pi/lawnberry/data
ExecStartPre=/bin/mkdir -p /home/pi/lawnberry/logs  
ExecStartPre=/bin/mkdir -p /home/pi/lawnberry/config
PermissionsStartOnly=true
ExecStartPre=/bin/chown -R pi:pi /home/pi/lawnberry

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/pi/lawnberry

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-database

[Install]
WantedBy=multi-user.target