# src/safety/lawnberry-safety.service
[Unit]
Description=Lawnberry Safety Monitoring System
Documentation=https://github.com/lawnberry/safety-system
After=network-online.target mosquitto.service
Wants=network-online.target
Requires=mosquitto.service
StartLimitIntervalSec=60

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/var/lib/lawnberry
# Use writable CWD to avoid lgpio .lgd-nfy* pipe issues if GPIO is touched
Environment=GPIOZERO_PIN_FACTORY=lgpio
Environment=XDG_RUNTIME_DIR=/run/lawnberry
Environment=LG_GPIO_CHIP=0
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.safety.safety_main
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StartLimitBurst=3

# Resource limits
MemoryMax=1024M
CPUQuota=50%
TasksMax=100

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry /run /tmp
RuntimeDirectory=lawnberry
RuntimeDirectoryMode=0755
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Performance settings
Nice=-5
IOSchedulingClass=1
IOSchedulingPriority=4

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-safety

[Install]
WantedBy=multi-user.target
Also=mosquitto.service
