# src/system_integration/lawnberry-auto-redeploy.service
[Unit]
Description=LawnBerryPi Auto Rebuild & Deploy Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
# Run from canonical runtime but watch the editable workspace via WATCH_ROOT
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=WATCH_ROOT=/home/pi/lawnberry
Environment=FAST_DEPLOY_SERVICE_TIMEOUT=10
Environment=FAST_DEPLOY_MAX_SECONDS=70
ExecStart=/bin/bash -lc '/bin/bash /opt/lawnberry/scripts/auto_rebuild_deploy.sh'
Restart=always
RestartSec=5
TimeoutStartSec=20
TimeoutStopSec=20
NoNewPrivileges=false
ProtectSystem=full
# Allow read access to home for watching; write is not needed
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID
ProtectKernelTunables=true
ProtectControlGroups=true
LockPersonality=true
MemoryMax=300M
ReadWritePaths=/opt/lawnberry /etc/systemd/system /var/log/lawnberry

# Additional Bookworm Security Features
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
# Allow sudo to elevate for service restarts; RestrictSUIDSGID would force no_new_privs and break sudo
RestrictSUIDSGID=false
SystemCallArchitectures=native
UMask=0027

[Install]
WantedBy=multi-user.target
