# src/system_integration/lawnberry-system.service
[Unit]
Description=Lawnberry System Integration Manager
Documentation=https://github.com/lawnberry/system-integration
After=network-online.target dbus.service multi-user.target
Wants=network-online.target dbus.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
EnvironmentFile=-/opt/lawnberry/.env
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.system_integration.system_manager
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitIntervalSec=300

# Resource limits
MemoryMax=512M
CPUQuota=100%
TasksMax=200

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry /opt/lawnberry/config /opt/lawnberry/data
# If your app ever needs to write to /run/lawnberry or other locations, add them here.
# Security settings
## NOTE: Runs as user pi from a system unit. If your code ever tries to reach the user session bus (systemctl --user ...),
## there is often no user bus for a background system unit. System bus is fine; user bus needs an active user session/linger.
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Performance settings - Bookworm optimized
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=2
CPUAccounting=true
CPUWeight=200

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-system

# Startup/shutdown timeouts
TimeoutStartSec=90
TimeoutStopSec=30

# Bookworm-specific optimizations
PrivateUsers=false
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0027

[Install]
WantedBy=multi-user.target
