# src/sensor_fusion/lawnberry-sensor-fusion.service
[Unit]
Description=Lawnberry Sensor Fusion Engine
Documentation=https://github.com/lawnberry/sensor-fusion
After=lawnberry-hardware.service lawnberry-communication.service
Requires=lawnberry-hardware.service lawnberry-communication.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.sensor_fusion.fusion_engine
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StartLimitBurst=3

# Resource limits
MemoryMax=384M
CPUQuota=100%
TasksMax=100

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Performance settings
Nice=-2
IOSchedulingClass=1
IOSchedulingPriority=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-sensor-fusion

# Watchdog removed (sd_notify not used by service)

[Install]
WantedBy=multi-user.target
