# web-ui/lawnberry-api.service
[Unit]
Description=Lawnberry Web API Backend
After=network.target mosquitto.service redis.service
Wants=mosquitto.service redis.service
StartLimitIntervalSec=0

[Service]
Type=exec
Restart=always
RestartSec=5
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
ExecStart=/opt/lawnberry/venv/bin/python /opt/lawnberry/src/web_api/run_server.py
Environment=PYTHONPATH=/opt/lawnberry
Environment=LOG_LEVEL=INFO
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=DEBUG=True

# Allow access to home directory for development environment
ProtectHome=false

# Security settings (relaxed for development)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/lawnberry /tmp /opt/lawnberry
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Device access for camera/video (API streams MJPEG)
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/video1 rw
DeviceAllow=/dev/video2 rw
DeviceAllow=/dev/video3 rw
DeviceAllow=/dev/media0 rw
DeviceAllow=/dev/media1 rw

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-api

[Install]
WantedBy=multi-user.target
