# src/hardware/lawnberry-sensor.service
[Unit]
Description=Lawnberry Hardware Sensor Service
After=network-online.target mosquitto.service
Wants=network-online.target mosquitto.service
StartLimitIntervalSec=0

[Service]
Type=exec
User=pi
Group=pi
# Use a writable working directory so lgpio can create its notification pipes (.lgd-nfy*)
WorkingDirectory=/var/lib/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
Environment=LOG_LEVEL=INFO
# Ensure gpiozero uses lgpio on Bookworm (more reliable on Pi 4/5)
Environment=GPIOZERO_PIN_FACTORY=lgpio
Environment=LAWNBERY_TOF_NO_GPIO=always
# Optional: enable one-time address assignment if a sensor is missing at startup (0/1)
Environment=LAWNBERY_TOF_AUTO_ASSIGN_ON_MISSING=1
# Optional: max seconds to allow the auto-assign script to run
Environment=LAWNBERY_TOF_AUTO_ASSIGN_TIMEOUT_S=60
ExecStart=/opt/lawnberry/venv/bin/python -m src.hardware.sensor_service
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StartLimitBurst=5
TimeoutStopSec=20
KillMode=mixed

# Resource limits
MemoryMax=256M
CPUQuota=100%
TasksMax=100

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# Allow writes to logs, state dir (working dir), and device nodes
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry /dev/i2c-1 /dev/ttyAMA0 /dev/ttyAMA4 /dev/ttyACM0 /dev/ttyACM1
DeviceAllow=/dev/i2c-1 rw
DeviceAllow=/dev/ttyAMA0 rw
DeviceAllow=/dev/ttyAMA4 rw
DeviceAllow=/dev/ttyACM0 rw
DeviceAllow=/dev/ttyACM1 rw
DeviceAllow=/dev/gpiochip0 rw
DeviceAllow=/dev/gpiochip1 rw
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Hardware access
SupplementaryGroups=gpio i2c spi dialout

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-sensor

[Install]
WantedBy=multi-user.target
