# src/hardware/lawnberry-hardware.service
[Unit]
Description=Lawnberry Hardware Interface Layer
Documentation=https://github.com/lawnberry/hardware-interface
After=network.target lawnberry-communication.service lawnberry-data.service
Requires=lawnberry-communication.service lawnberry-data.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.hardware.hardware_interface
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
StartLimitBurst=3

# Resource limits
MemoryMax=256M
CPUQuota=75%
TasksMax=100

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry /dev
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Hardware access
SupplementaryGroups=gpio i2c spi dialout

# Performance settings
Nice=-3
IOSchedulingClass=1
IOSchedulingPriority=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-hardware

# Watchdog
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
