# src/hardware/lawnberry-hardware.service
[Unit]
Description=LawnBerry Hardware Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lawnberry
Group=lawnberry
# lgpio creates .lgd-nfy* in the CWD; use a writable location
WorkingDirectory=/var/lib/lawnberry
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/lawnberry
Environment=GPIOZERO_PIN_FACTORY=lgpio
# Use the sensor service entrypoint which is the long-running hardware publisher
ExecStart=/opt/lawnberry/venv/bin/python -m src.hardware.sensor_service
Restart=on-failure
RestartSec=3
SupplementaryGroups=gpio i2c spi dialout
# Hardening options
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
ReadOnlyPaths=/
ReadWritePaths=/var/lib/lawnberry /var/log/lawnberry /run /tmp
DevicePolicy=closed
DeviceAllow=/dev/i2c-1 rw
DeviceAllow=/dev/gpiomem rw
DeviceAllow=/dev/gpiochip* rw
DeviceAllow=/dev/tty* rw
MemoryMax=512M

[Install]
WantedBy=multi-user.target
 
