# src/hardware/lawnberry-hardware.service
{# Updated service to use writable working directory and correct entrypoint #}
[Unit]
Description=LawnBerry Hardware Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lawnberry
Group=lawnberry
# lgpio creates .lgd-nfy* in the CWD; use a writable location
WorkingDirectory=/var/lib/lawnberry
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/lawnberry
Environment=GPIOZERO_PIN_FACTORY=lgpio
# Use the sensor service entrypoint which is the long-running hardware publisher
ExecStart=/opt/lawnberry/venv/bin/python -m src.hardware.sensor_service
Restart=on-failure
RestartSec=3
# Hardening options
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
ReadOnlyPaths=/
ReadWritePaths=/var/lib/lawnberry /var/log/lawnberry /run /tmp
DeviceAllow=/dev/i2c-1 rw
DeviceAllow=/dev/gpiomem rw
DeviceAllow=/dev/tty* rw
MemoryMax=512M

[Install]
WantedBy=multi-user.target
SupplementaryGroups=gpio i2c spi dialout

# Performance settings
Nice=-3
IOSchedulingClass=1
IOSchedulingPriority=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-hardware

# Watchdog
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
