# src/power_management/lawnberry-power.service
[Unit]
Description=Lawnberry Power Management System
Documentation=https://github.com/lawnberry/power-management
After=lawnberry-hardware.service lawnberry-communication.service
Requires=lawnberry-hardware.service lawnberry-communication.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/lawnberry
Environment=PYTHONPATH=/opt/lawnberry
Environment=LAWNBERRY_ENV=production
ExecStart=/opt/lawnberry/venv/bin/python -m src.power_management.power_service
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StartLimitBurst=3

# Resource limits
MemoryMax=256M
CPUQuota=50%
TasksMax=75

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/lawnberry /var/lib/lawnberry
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Performance settings
Nice=0
IOSchedulingClass=1
IOSchedulingPriority=4

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-power

# Watchdog removed (sd_notify not used by service)

[Install]
WantedBy=multi-user.target
