# src/web_api/lawnberry-api.service
[Unit]
Description=Lawnberry Web API Backend
After=network-online.target mosquitto.service
Wants=network-online.target mosquitto.service
StartLimitIntervalSec=0

[Service]
Type=exec
Restart=always
RestartSec=5
User=pi
Group=pi
WorkingDirectory=/home/pi/lawnberry_pi
ExecStart=/home/pi/lawnberry_pi/venv/bin/python /home/pi/lawnberry_pi/src/web_api/run_server.py
EnvironmentFile=-/home/pi/lawnberry_pi/.env
Environment=PYTHONPATH=/home/pi/lawnberry_pi/src
Environment=LOG_LEVEL=INFO
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=LAWNBERRY_ENV=production

# Security hardening (Bookworm-compatible)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/home/pi/lawnberry_pi/web-ui/dist /var/log/lawnberry /var/lib/lawnberry

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawnberry-api

[Install]
WantedBy=multi-user.target
